# src/systemd/rhosts.path
[Unit]
Description="Monitor the /etc/passwd file for changes"

[Path]
PathModified=/etc/rhosts/rhosts.cfg
Unit=rhosts.service

[Install]
WantedBy=multi-user.target
